# deploy/launchd/com.cmux-persist.watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cmx.watch</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{BINARY}}</string>
        <string>watch</string>
        <string>autosave</string>
        <string>--interval</string>
        <string>5m</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <dict>
        <key>PathState</key>
        <dict>
            <key>/tmp/cmux.sock</key>
            <true/>
        </dict>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/cmx-watch.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/cmx-watch.log</string>
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
